# e2e/com.clawvault.daemon.test.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.clawvault.daemon</string>
    <key>ProgramArguments</key>
    <array>
        <string>__DAEMON_BINARY__</string>
    </array>
    <key>MachServices</key>
    <dict>
        <key>com.clawvault.daemon</key>
        <true/>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/clawvault-e2e-daemon.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/clawvault-e2e-daemon.err</string>
</dict>
</plist>
